# supervisord.conf
[supervisord]
nodaemon=true



[program:glassfish]
command=bash -c ". /opt/glassfish4/glassfish/bin/startserv"
